# config/Code - OSS/User/History/-79b66d17/cqMO.service
[Unit]
Description=Display colour temperature adjustment
PartOf=graphical-session.target
After=graphical-session.target NetworkManager-wait-online.service

[Service]
ExecStart=%h/.local/bin/gammastep-wrapper.sh

[Install]
WantedBy=graphical-session.target